# rclone/Library/LaunchAgents/org.rclone.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>KeepAlive</key>
    <false/>
    <key>Label</key>
    <string>org.rclone</string>
    <key>ProgramArguments</key>
    <array>
      <string>/usr/local/bin/rclone</string>
      <string>mount</string>
      <string>--checksum</string>
      <string>--vfs-cache-mode</string>
      <string>writes</string>
      <string>aws:yuanl-org</string>
      <string>/Users/liyuanliang/org/</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
  </dict>
</plist>
